# retrocard.umount.service
[Unit]
Description=Relaunch emulationstation after retrocard usb keys unmount
Requires=media-retrocard.umount
After=media-retrocard.umount

[Service]
ExecStart=/usr/bin/emulationstation

[Install]
WantedBy=media-retrocard.umount
