# retrocard.mount.service
[Unit]
Description=Mount retrocard usb keys on RETROCARD folder
Requires=target.mount
Before=target.mount

[Install]
WantedBy=target.mount

[Mount]
What=/dev/disk/by-label/RETROCARD
Where=/media/RETROCARD
