# retrocard.umount.service
[Unit]
Description=Relaunch emulationstation after retrocard usb keys unmount
Requires=media-usb0.umount
After=media-usb0.umount

[Service]
ExecStart=/usr/bin/emulationstation

[Install]
WantedBy=media-usb0.umount
