# retrocard.launch.service
[Unit]
Description=Launch from retrocard
Requires=mount.target
After=mount.target

[Service]
ExecStart=/home/pi/retrocard/launch.sh

[Install]
WantedBy=mount.target
